# systemd/jf-manager-backup.timer
[Unit]
Description=JF-Manager Daily Backup Timer

[Timer]
# Run daily at 2 AM
OnCalendar=daily
OnCalendar=*-*-* 02:00:00
Persistent=true

# Randomize start time by up to 15 minutes to avoid load spikes
RandomizedDelaySec=15min

[Install]
WantedBy=timers.target
